# crontab
0 8 * * * /var/www/comments.thatsaspicymeatball.com/generate.rb >> /dev/null 2>&1
